# fs/etc/services.d/configurer/run
#!/usr/bin/execlineb -P

with-contenv

python3 /usr/lib/gluster-init/gluster-init.py -vvv