# fs/etc/services.d/glusterd/run
#!/usr/bin/execlineb -P

exec glusterd --no-daemon